# the-luminal-archive/memory/ACTIVE_SYSTEM/gritz-memory-llm.service
[Unit]
Description=Gritz LLM Memory Service - AI-Enhanced Memory Processing
After=network.target

[Service]
Type=simple
User=ubuntumain
WorkingDirectory=/home/ubuntumain/Documents/Github/project-sanctuary/the-luminal-archive/memory/ACTIVE_SYSTEM
ExecStart=/home/ubuntumain/Documents/Github/project-sanctuary/the-luminal-archive/memory/llm_venv/bin/python /home/ubuntumain/Documents/Github/project-sanctuary/the-luminal-archive/memory/ACTIVE_SYSTEM/llm_service.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target